# rudder-contrib-backup/etc/cron.d/rudder-contrib-backup
01 02 * * *	root /usr/sbin/rudder-contrib-backup 2>&1 | logger -p daemon.info -t rudder-contrib-backup
